# fsoaudiod/data/fsoaudiod.service
[Unit]
Description=FSO daemon implementing the audio API
After=dbus.socket

[Service]
ExecStart=/usr/sbin/fsoaudiod

[Install]
WantedBy=basic.target
