# fsosystemd/data/fsosystemd.service
[Unit]
Description=FSO daemon implementing the basic system API

[Service]
ExecStart=/usr/sbin/fsosystemd

[Install]
WantedBy=basic.target
